# battle_of_the_dragons_revived/app/supervisord.conf
[supervisord]
nodaemon=true

[program:nwserver-linux]
command=/nwserver/bin/linux-x86/nwserver-linux
autostart=true
autorestart=true
